# workloads/vla-serving-host/ansible/templates/openvla-server.container.j2
# Podman Quadlet unit — systemd-managed podman container.
# Generated by Ansible from templates/openvla-server.container.j2.
# Rendered location: /etc/containers/systemd/openvla-server.container
[Unit]
Description=Host-native OpenVLA-7B serving runtime (ADR-026)
Wants=network-online.target
After=network-online.target

[Container]
ContainerName=openvla-server
Image={{ container_image }}
# Host networking so the SNO CNI range can reach the listen port
# through the podman bridge without additional NAT rules.
Network=host
# ROCm GPU + NPU device passthrough.
AddDevice=/dev/kfd
AddDevice=/dev/dri
AddDevice=/dev/accel
# Model cache persists between restarts.
Volume={{ model_cache }}:/var/cache/vla-models:Z
User={{ runtime_uid }}
Group={{ runtime_uid }}
Environment=HF_HOME=/var/cache/vla-models/huggingface
Environment=SERVICE_NAME=openvla-server
Environment=LOG_LEVEL=INFO
Environment=PORT={{ listen_port }}
Environment=OPENVLA_WEIGHTS=openvla/openvla-7b
# ROCm tuning — Strix Halo (gfx1151) explicitly targeted.
Environment=HSA_OVERRIDE_GFX_VERSION=11.5.1
Environment=PYTORCH_ROCM_ARCH=gfx1151
Environment=ROCR_VISIBLE_DEVICES=0

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target default.target
